# deployment/raspberry-pi/roles/moon_service_host/templates/moon-service-docker-firewall.service.j2
[Unit]
Description=Block Moon Service containers from IP-authorized SMTP relays
Documentation=https://github.com/rapucha/moon-service/blob/main/deployment/raspberry-pi/README.md
BindsTo=docker.service
After=docker.service
PartOf=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/moon-service-docker-firewall
RemainAfterExit=true
TimeoutStartSec=30s
UMask=0077
RuntimeDirectory=moon-service
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
NoNewPrivileges=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadWritePaths=/run/moon-service -/run/xtables.lock
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
