# deployment/raspberry-pi/roles/moon_service_host/templates/moon-service-public-notify.timer.j2
[Unit]
Description=Retry queued Moon Service operator notifications

[Timer]
OnActiveSec=2min
OnUnitInactiveSec=5min
AccuracySec=30s
Unit=moon-service-public-notify.service

[Install]
WantedBy=timers.target
